# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log ; (default log path)
loglevel=info                    ; log level
pidfile=/var/run/supervisord.pid ; (default pid file)

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; (the path to the socket file)

[program:app]
;command=flask run --host=0.0.0.0
command=python -m app
directory=/workspace
autostart=true
autorestart=true
stderr_logfile=/workspace/logs/app.err.log
stdout_logfile=/workspace/logs/app.out.log